# debian/ffwd.service.in
[Unit]
Description=fast-forward monitoring agent
Documentation=https://github.com/spotify/ffwd

[Service]
ExecStart=/usr/bin/ffwd-RELEASE /etc/ffwd/ffwd-RELEASE.yaml
ExecStopPost=/bin/sleep 10
Restart=on-failure

[Install]
WantedBy=multi-user.target
